# docker/supervisor-lemp.conf
[program:nginx]
command=nginx

[program:php-fcgi]
command=php-fcgi

[program:mysqld]
command=mysqld_safe